# packages/rocknix/system.d/rocknix-battery-logging.service
[Unit]
Description=Log battery statistics

[Service]
Type=simple
ExecStart=/usr/bin/rocknix-battery-log

[Install]
WantedBy=multi-user.target
